# etc/supervisor/conf.d/dshd.conf
[program:SSH]
command=/usr/sbin/sshd -D 
directory=/root/
