# etc/supervisor/conf.d/dshd.conf
[program:SSH]
command=bash -c '/usr/sbin/sshd -D'
directory=/root/
