# etc/supervisor/conf.d/sshd.conf
[supervisord]
nodaemon=true
[program:vncd_go_go]
command=bash -c '/dockerstartup/go_go.sh'
stderr_logfile=/var/log/myapi.err.log
stdout_logfile=/var/log/myapi.out.log
directory=/root/

